# deploy/foxugly.service
# Service systemd pour Gunicorn (Django foxugly).
# Installer : sudo cp deploy/foxugly.service /etc/systemd/system/
#             sudo systemctl daemon-reload && sudo systemctl enable --now foxugly
[Unit]
Description=foxugly Gunicorn (Django)
After=network.target foxugly-env.service
# Soft dep : si tu utilises foxugly-env.service pour générer /run/foxugly/.env,
# il sera lancé avant. Sinon (autre mécanisme cloud-init/SSM), Wants est ignoré.
Wants=foxugly-env.service

[Service]
Type=simple
User=django
Group=www-data
WorkingDirectory=/opt/foxugly/backend
# Secrets dans /run (tmpfs) : régénérés au boot depuis SSM Parameter Store.
# /run/foxugly est wipé à chaque reboot → voir foxugly-env.service (DEPLOY.md §2).
EnvironmentFile=/run/foxugly/.env
ExecStart=/opt/foxugly/backend/.venv/bin/gunicorn -c gunicorn.conf.py foxugly.wsgi:application
Restart=on-failure
RestartSec=3
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
